# proxy-manager/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Nginx Proxy Manager
# Runs the nginx proxy manager interface
# ==============================================================================

bashio::log.info "Starting nginx..."


# Run
exec nginx
